# jbot-kr-dns-vote-jbot-kr/deploy/jbot-vote.service
[Unit]
Description=JBOT School Vote
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/jbot-vote
Environment=PORT=3000
Environment=ADMIN_PASSWORD=change-this-admin-password
ExecStart=/usr/bin/node /opt/jbot-vote/server.js
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
